# conf.dev/supervisor-beat.conf
[program:alfa-lombard-beat]
directory = /usr/local/flexites/projects_ext/alfa-omega-lombard/server/proj
command = /usr/local/flexites/projects_ext/alfa-omega-lombard/server/env/bin/celery -A conf beat -l INFO --pidfile="/usr/local/flexites/projects_ext/alfa-omega-lombard/var/run/celerybeat.pid" --schedule="/usr/local/flexites/projects_ext/alfa-omega-lombard/var/run/celerybeat.db"
user = www-data
stdout_logfile = /usr/local/flexites/projects_ext/alfa-omega-lombard/var/log/beat.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 5
stderr_logfile = /usr/local/flexites/projects_ext/alfa-omega-lombard/var/log/beat.error.log
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups = 5
startsecs = 10
autostart = true
autorestart = true
priority = 997
stopsignal = QUIT
